# bootstrap/ansible/roles/k3s-snapshot-sync/templates/k3s-snapshot-sync.service.j2
[Unit]
Description=Sync k3s etcd snapshots to jaguar
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/k3s-snapshot-sync
StandardOutput=journal
StandardError=journal
